# dockerfiles/scrape.crontab
0 */2 * * * cd /app && poetry run scrape
